# etc/udev/rules.d/50-wol.rules
ACTION=="add|change", SUBSYSTEM=="net", KERNEL=="eth*", RUN+="/sbin/ethtool -s %k wol d"
