# etc/udev/rules.d/95-monitor-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", ENV{HOTPLUG}=="1", RUN+="/home/xlz/bin/drm.sh"
